# udev/rules.d/50-flash.rules
KERNEL=="sd[a-z][0-9]", ATTRS{vendor}=="Kingston", RUN+="/bin/mount /dev/%k /mnt/.ssh"
